# .config/udev-rules/99-battery.rules
SUBSYSTEM=="power_supply", ACTION=="change", RUN+="/home/louis/.local/bin/update-battery" 
