# .config/udev-rules/99-connection.rules
#SUBSYSTEM=="net", RUN+="/home/louis/.local/bin/update-connection"

