# .config/udev-rules/99-display.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="/home/louis/.local/bin/update-screen"
